# deploy/systemd/ebpffls.service
[Unit]
Description=ebpffls eBPF anti-ransomware guard
Documentation=https://github.com/panzhenyu12/ebpffls
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ebpffls monitor --config /etc/ebpffls/ransomware.yaml --dry-run=false
Restart=always
RestartSec=5s
KillMode=process
LimitMEMLOCK=infinity
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/run /var/log /var/lib/ebpffls
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_SYS_ADMIN CAP_SYS_RESOURCE CAP_SYS_PTRACE CAP_KILL CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_ADMIN CAP_SYS_RESOURCE CAP_SYS_PTRACE CAP_KILL CAP_DAC_READ_SEARCH

[Install]
WantedBy=multi-user.target
